# contrib/chaff-hygiene.service
[Unit]
Description=chaff git hygiene pass

[Service]
Type=oneshot
ExecStart=%h/.claude/scripts/chaff-cron.sh
StandardOutput=journal
StandardError=journal
